# .config/systemd/user/drive.service
[Unit]
Description=Run Rclone Service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/home/user/drive.sh
