# .config/systemd/user/drive.service
[Unit]
Description=Run Rclone Service

[Service]
Type=simple
ExecStart=/bin/bash /usr/local/bin/drive
