# .config/systemd/user/drive.service
[Unit]
Description=Drive Mount Service
After=i3-session.target

[Service]
Type=oneshot
ExecStart=sh /home/user/drive.sh
RemainAfterExit=true
