# .config/systemd/user/drive.service
[Unit]
Description=Run Rclone Service

[Service]
ExecStart=/home/user/drive.sh
