# .config/systemd/user/drive.service
[Unit]
Description=Google Drive Mount Service

[Service]
Type=oneshot
ExecStart=sh ${HOME}/drive.sh
