# .config/systemd/user/feh.service
[Unit]
Description=feh
After=i3-session.target
Requires=i3-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/feh --bg-scale ${HOME}/.config/background.jpg
#ExecReload=/usr/bin/feh --bg-scale ${HOME}/.config/background.jpg
#Restart=always

ProtectSystem=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
