# .config/systemd/user/feh.service
[Unit]
Description=feh
Requires=i3-session.target
After=i3-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/feh --bg-fill /home/user/.config/background.jpg
ExecReload=/usr/bin/feh --bg-fill /home/user/.config/background.jpg
ExecStop=/usr/bin/feh --bg-fill /home/user/.config/background.jpg
RemainAfterExit=true

ProtectSystem=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
