# .config/systemd/user/grive.service
[Unit]
Description=Grive Service

[Service]
Type=oneshot
ExecStart=/usr/bin/grive -p ${HOME}/Drive
