# .config/systemd/user/i3-session.target
[Unit]
Description=i3
After=default.target
BindsTo=graphical-session.target
#Wants=feh.service
#Wants=graphical-session-pre.target
Wants=picom.service
#Wants=restart-service@dunst.path
#Wants=xss-lock.service
