# .config/systemd/user/watch.timer
[Unit]
Description=15 seconds timer to sync

[Timer]
OnBootSec=2min
OnUnitInactiveSync=15sec
Unit=grive.service

[Install]
WantedBy=basic.target
