# config/rootfs/etc/cron.d/grive
* * * * * umut cd ${HOME}/Drive && grive
